# templates/systemd.service.j2
#
# {{ ansible_managed }}
#
# This file was generated by Ansible for {{ ansible_fqdn }}
# Do NOT modify this file by hand!
#
# The uWSGI emperor systemd configuration
#

[Unit]
Description=uWSGI Emperor
After=syslog.target

[Service]
#ExecStart=/usr/local/bin/uwsgi --ini /etc/uwsgi-emperor/emperor.ini
ExecStart=/usr/bin/uwsgi --ini /etc/uwsgi-emperor/emperor.ini
Restart=always
KillSignal=SIGQUIT
Type=notify
#StandardError=syslog
NotifyAccess=all
#RuntimeDirectory=uwsgi
User=root
Group=root

[Install]
WantedBy=multi-user.target

